# cmd/upstart/gomon.conf
description "Pluggable Go Server and Service Monitoring"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5

exec /usr/local/bin/gomon --config /usr/local/etc/gomon.gcfg
